# nominatim-ac.conf
description "Nominatim AC"
author "Tony McCrae <tony@eelpieconsulting.co.uk>"

start on net-device-up

respawn

exec java -Xmx1g -jar /opt/nominatim-ac/nominatim-ac.jar --server.port=10010
